# package/base-files/files/etc/init.d/telnetd
#!/bin/sh /etc/rc.common
START=99
start() {
        telnetd -l /bin/login -p 23
}
stop() {
        killall telnetd 2>/dev/null
}
